# systemd/soundmask.service
[Unit]
Description=SoundMask Calendar Sound Masking
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=soundmask
Group=soundmask
WorkingDirectory=/opt/SoundMask
EnvironmentFile=-/etc/soundmask/soundmask.env
ExecStart=/opt/SoundMask/.venv/bin/soundmask serve
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
StandardOutput=append:/var/lib/soundmask/logs/service.log
StandardError=append:/var/lib/soundmask/logs/service.log
Restart=always
RestartSec=5
PrivateTmp=true
StateDirectory=soundmask

[Install]
WantedBy=multi-user.target
