# systemd/soundmask-update-check.service
[Unit]
Description=SoundMask update check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/SoundMask
EnvironmentFile=-/etc/soundmask/soundmask.env
ExecStart=/opt/SoundMask/.venv/bin/python -m app.update_manager check
