# systemd/soundmask-update-install.service
[Unit]
Description=SoundMask update install
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/SoundMask
EnvironmentFile=-/etc/soundmask/soundmask.env
ExecStart=/opt/SoundMask/.venv/bin/python -m app.update_manager install
StandardOutput=append:/var/lib/soundmask/logs/updates.log
StandardError=append:/var/lib/soundmask/logs/updates.log
